# tools/archive_reservations.service
[Unit]
Description=Archivage des anciennes reservations Vehicules

[Service]
Type=oneshot
WorkingDirectory=/workspace/vehicules
ExecStart=/workspace/vehicules/venv/bin/python tools/archive_reservations.py --days 7
