# tools/backup_db.service
[Unit]
Description=Sauvegarde de la base Vehicules

[Service]
Type=oneshot
WorkingDirectory=/path/to/vehicules
Environment=REMOTE_URI=gdrive:vehicules-backups
ExecStart=/usr/bin/env bash tools/backup_db.sh
